# rrd-api.service
[Unit]
Description=RRD Web API Server
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/rrd-api
Environment=RRD_BASE_DIR=/var/lib/rrd
ExecStart=/usr/local/bin/uvicorn rrd_api_server:app --host 0.0.0.0 --port 8080
Restart=always

[Install]
WantedBy=multi-user.target
